# systemd/rtx.service
[Unit]
Description=homelab-rtx GPU telemetry service

[Service]
Type=simple
User=rtx
Group=rtx
WorkingDirectory=/srv/rtx
EnvironmentFile=-/etc/rtx/rtx.env
ExecStart=/srv/rtx/scripts/up.sh
Restart=on-failure
RestartSec=5s
KillSignal=SIGTERM
TimeoutStopSec=30s
StandardOutput=journal
StandardError=journal
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
